# S99drcom
#!/bin/sh /etc/rc.common
# put it into /etc/rc.d
START=99
APP="python /usr/bin/drcom.py" 
start() {
	sleep 30
	start-stop-daemon -S -x $APP
}

stop() {
	killall python
}
